# upstart/gohub.conf
description "GitHub webhook deployment server"
author "Matt Harzewski"

start on startup
stop on shutdown
respawn

script
        chdir /etc/gohub
        exec gohub --port 8392
end script